# mirror-selector@.service
[Unit]
Description=change apt mirror

[Service]
Slice=session.slice
Type=notify
ExecStart=/usr/bin/mirror-selector-commander -a %i
